# etc/crontab-tiger.cron
SHELL=/bin/sh
PATH=/sbin:/bin:/usr/local/sbin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin:/root/bin
# Nightly backup script
0 01 * * * ( backup --config-file /usr/local/etc/backup-tiger.conf --level `date +\%w` )
